# crontab
*/5 * * * * /usr/bin/bash /home/paul/personal/earthquake-data/get-eq-data.bash
